# linux/flatpak/.config/autostart/com.bitwarden.desktop.desktop
[Desktop Entry]
Exec=flatpak run com.bitwarden.desktop
Name=com.bitwarden.desktop
Type=Application
X-Flatpak=com.bitwarden.desktop
